# telegram-monitor.service
[Unit]
Description=Telegram Price Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/telegram-price-monitor
EnvironmentFile=/home/ubuntu/telegram-price-monitor/.env
ExecStart=/home/ubuntu/telegram-price-monitor/venv/bin/python monitor.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
